# rootfs/etc/s6-overlay/s6-rc.d/znc/run
#!/usr/bin/with-contenv sh

# change chdir
cd /config

# run (from /usr/local/bin) with limited permissions
exec s6-setuidgid $PUID:$PGID znc --foreground
